# package/fon-luci/luci/i18n/portuguese/root/etc/uci-defaults/luci-i18n-portuguese
#!/bin/sh
uci batch <<-EOF
	set luci.languages.pt=Portugese
        commit luci
EOF
	
