# package/fonstate/files/etc/init.d/fonstate_early
#!/bin/sh /etc/rc.common

START=42

start() {
		echo -n "Offline" > /tmp/run/fonstate
		echo -n "Offline" > /tmp/run/hotspot
}
