# package/fon-luci/luci/modules/admin-fon/root/etc/uci-defaults/zzz-upgrade-migrate
#!/bin/sh

# This script calls the configuration migration script, to perform any
# post-upgrade config migrations. It is named such that it runs _after_
# the zzz-upgrade script that is generated by save-config.sh during a
# an upgrade before the reboot.

/sbin/migrate-config.sh
