# package/fonbase/files/etc/hotplug.d/iface/20-hotspot
[ "$ACTION" == "ifup" ] && {
		[ "$INTERFACE" = "hotspot" ] && {
			iptables -t mangle -I PREROUTING 1 -i ath0 -j DROP
		}
}
