# package/fon-luci/luci/themes/fon/root/etc/uci-defaults/luci-theme-fon
#!/bin/sh
uci batch <<-EOF
	set luci.themes.FON=/luci-static/fon
        commit luci
EOF
	
