# package/fon-luci/luci/i18n/basque/root/etc/uci-defaults/luci-i18n-basque
#!/bin/sh
uci batch <<-EOF
	set luci.languages.eu=Euskara
        commit luci
EOF
	
