# package/fon-luci/luci/i18n/dutch/root/etc/uci-defaults/luci-i18n-dutch
#!/bin/sh
uci batch <<-EOF
	set luci.languages.nl=Dutch
        commit luci
EOF
	
