# package/fon-luci/luci/applications/luci-wizard-fonera2/root/etc/uci-defaults/fon2-wizard
uci set wizard.fon2wiz=wizard
uci set wizard.fon2wiz.callonce=step1
uci commit wizard
