# package/fon-luci/luci/i18n/zh_TW/root/etc/uci-defaults/luci-i18n-zh_tw
#!/bin/sh
uci batch <<-EOF
	set luci.languages.zh="正體中文"
        commit luci
EOF
	
