# package/fonbase/files/etc/hotplug.d/net/10-hotspot
case "$ACTION" in
	add|register)
		[ "$INTERFACE" = "tun0" ] && {
			/usr/lib/qos/generate.sh interface hotspot | sh
			iptables -t mangle -I PREROUTING 1 -i ath0 -j DROP
		}
		;;
esac
