# package/fonbase/files/etc/hotplug.d/net/90-mgmtnetwork
[ "$INTERFACE" = "br-lan" ] && ifconfig br-lan 2>/dev/null >/dev/null && {
/etc/init.d/mgmtnetwork start
}
