# package/fonbase/files/etc/uci-defaults/host
echo "192.168.10.1 registerlafonera.fon.com www.registerlafonera.fon.com fonera lafonera myfonera" >> /etc/hosts
