# package/fonbase/files/etc/uci-defaults/wpa
# this gets the wpa key based on serial
/lib/fon/config.sh fon
/lib/fon/config.sh wireless
uci commit
